# tcp_server.service
[Unit]
Description=Custom TCP Server
After=network.target

[Service]
ExecStart=/usr/bin/python3 tcp_server.py
WorkingDirectory=/path/to/your/script_directory
Restart=always
User=your_username
Group=your_groupname
StandardOutput=append:/var/log/tcp_server.log
StandardError=append:/var/log/tcp_server_error.log
	
[Install]
WantedBy=multi-user.target
